# config/systemd/user/trackpoint-adjuster.service
[Unit]
Description=Trackpoint Speed Adjuster

[Service]
ExecStart=/opt/trackpoint-adjuster/apply.sh &

[Install]
WantedBy=default.target
